# deploy/lightsail/polis-website-full.service
[Unit]
Description=Polis Full Website Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/polis-website/backend
EnvironmentFile=/etc/polis-website/backend.env
ExecStart=/usr/bin/node /var/www/polis-website/backend/src/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
